# mypcsetup/550flatpak/sysconfig/flatpak-update.service
[Unit]
Description=Flatpak update
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update -y
